# configs/worker/worker.ini
[program:worker]
user=tine20
command=php -d include_path=/etc/tine20/:/usr/share/tine20 /usr/share/tine20/worker.php --config /etc/tine20/conf.d/actionQueue.ini
autostart=true
autorestart=true
priority=15